# backend/companion-relay.service
# Install to /etc/systemd/system/companion-relay.service then:
#   systemctl daemon-reload && systemctl enable --now companion-relay
# Adjust WorkingDirectory / paths if you didn't use /root/companion-relay.
[Unit]
Description=companion relay backend (channel <-> PWA)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/companion-relay
EnvironmentFile=/root/companion-relay/relay.env
ExecStart=/root/companion-relay/venv/bin/python -m uvicorn app:app --host 127.0.0.1 --port 3011 --timeout-graceful-shutdown 3
Restart=always
TimeoutStopSec=10
RestartSec=3

[Install]
WantedBy=multi-user.target
